# tools/builder/assets/ffs_registrar.service
[Unit]
Description=Registers FFS gadget descriptors (writes to ep0)
After=setup-gadget.service
Requires=setup-gadget.service

[Service]
Type=simple
User=registrar
Group=registrar
Environment="BROKER_LOG=debug"
ExecStart=/usr/local/bin/ffs_registrar

# Restart policy: always restart to ensure resilience
Restart=always
RestartSec=1

# Disable restart rate limiting to ensure continuous operation
StartLimitIntervalSec=0

# Logging to journal for debugging
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
