# tools/builder/assets/tezsign.service
[Unit]
Description=Runs tezsign
After=attach-gadget.service
Requires=attach-gadget.service

[Service]
Type=simple
User=tezsign
Group=tezsign
Environment="DATA_STORE=/data/tezsign"
ExecStart=/app/tezsign

# Restart policy: always restart to ensure resilience
Restart=always
RestartSec=2

# Disable restart rate limiting to ensure continuous operation
# This ensures the service keeps restarting even after multiple failures
StartLimitIntervalSec=0

# Resource limits - conservative values for embedded devices (RPi Zero 2W, RPi4, Radxa Zero 3)
# These are upper limits, not reservations
LimitNOFILE=4096
LimitNPROC=256

# Logging to journal for debugging (can be filtered by unit)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
